# iot/wechat/files/etc/init.d/wechat
#!/bin/sh /etc/rc.common

START=99
STOP=99

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

run_global()
{
    local enable
    config_get_bool enable $1 enable
    if [ "$enable" -eq "1" ]; then
        echo $enable
        wechat& 
        echo "wechat enabled." >/dev/console  
    else
        echo $enable
        echo "wechat disabled." >/dev/console
        killall -kill wechat
    fi  
}

start() 
{
    config_load wechat
    config_foreach run_global global
}

stop()
{
    echo ...
}
